# nemo-fix.service
[Unit]
Description=Nemo Fix

[Service]
ExecStart=/bin/bash /etc/nemo-fix/init.sh

[Install]
WantedBy=multi-user.target
